# deploy/hermes-mcp-login.service
# systemd unit — install as a SYSTEM unit (/etc/systemd/system/), runs as the
# agent's user so it shares HERMES_HOME and the token store. Co-location with
# the agent is mandatory: both processes read/write the same
# ~/.hermes/mcp-tokens/. Modelled on the existing hermes-dashboard.service so
# Ansible can template both the same way (User, HERMES_HOME, XDG_RUNTIME_DIR,
# EnvironmentFile). Adjust the uid in user@<uid>.service / XDG_RUNTIME_DIR and
# the paths to your deployment.

[Unit]
Description=hermes-mcp-login — browser-triggered OAuth login for MCP servers
After=network-online.target user@1001.service
Wants=network-online.target

[Service]
Type=simple
User=hermes
Group=hermes
WorkingDirectory=/home/hermes/.hermes
Environment=HERMES_HOME=/home/hermes/.hermes
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/hermes
Environment=XDG_RUNTIME_DIR=/run/user/1001
# Service knobs (HERMES_MCP_LOGIN_*); see .env.example. For a central reverse
# proxy on another host, set HERMES_MCP_LOGIN_HOST=0.0.0.0 + firewall the port.
EnvironmentFile=/home/hermes/.hermes/hermes-mcp-login.env
# Run inside the agent's venv interpreter (provides hermes-agent + mcp).
ExecStart=/home/hermes/.hermes/hermes-agent/venv/bin/python -m hermes_mcp_login
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
